# services.d/igdaemon/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
/usr/bin/igdaemon --driver=libusb --log-level=1 --send-timeout=1000 --receive-timeout=1000 --no-daemon
